# research_pulse/launchd/com.xinting.research-pulse.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.xinting.research-pulse</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/xintinghu/.research_pulse/run_server.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/xintinghu</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/xintinghu/.research_pulse/server.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/xintinghu/.research_pulse/server.err.log</string>
</dict>
</plist>
